# daily_restart.service
[Unit]
Description=Daily Restart

[Service]
Type=oneshot
ExecStart=/sbin/shutdown -r now
