# loveletter_extraction.timer
[Unit]
Description=Run loveletter_extraction daily at 2:15 AM

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true

[Install]
WantedBy=timers.target
